# crontab
5,20,35,50 * * * * ~/ooi/ooiData.sh 
# once an hour, delete old segments (-depth is for good style)
2 * * * * find ~/ooi/segments/ -depth -mtime +5 -delete
