# crontab
5,20,35,50 * * * * cd ~/ooi; ./ooiData.sh &> ./ooiData.log
